# eww/daemons/calendar-fetcher.timer
[Unit]
Description=Run Google Calendar Fetcher every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
